# deploy/systemd/simpler-objects-object-server@.service
[Unit]
Description=Simpler Objects object server (%i)
Documentation=https://github.com/ctengel/simpler-objects
After=network-online.target
Wants=network-online.target
# Mount-path convention: each instance's storage must be at /srv/simpler-objects/<instance>.
# RequiresMountsFor is evaluated at unit-load time and cannot expand EnvironmentFile
# variables; the convention is what lets the dependency be static. Override via
# a drop-in (`systemctl edit simpler-objects-object-server@%i`) if your mount
# lives elsewhere.
RequiresMountsFor=/srv/simpler-objects/%i

[Service]
Type=exec
EnvironmentFile=/etc/simpler-objects/object-server-%i.env
User=simpler-objects
Group=simpler-objects

# Defaults; values in the env file override these.
Environment=HOST=0.0.0.0
Environment=PORT=29171
Environment=WORKERS=1

# Refuse to start if a previous crash left orphan partial files or a garbled
# checksum file. Dry-run only — operator must run scrub with --delete-victims
# and/or --repair-checksums to clean up, then `systemctl reset-failed`.
ExecStartPre=/opt/simpler-objects/venv/bin/python -m simpler_objects.scrub ${OBJECT_DIRECTORY}

# Adjust the venv path to wherever you installed the package.
ExecStart=/opt/simpler-objects/venv/bin/uvicorn simpler_objects.object_server:app --host ${HOST} --port ${PORT} --workers ${WORKERS}

Restart=on-failure
RestartSec=5s

# Sandboxing
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
ReadWritePaths=/srv/simpler-objects/%i

[Install]
WantedBy=multi-user.target
